# deploy/hikari-bot.service
[Unit]
Description=HIKARI_BOT — NoneBot2 QQ Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/HIKARI_BOT
Environment=ENVIRONMENT=prod
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/root/.local/bin:/root/.cargo/bin
ExecStart=/root/.local/bin/uv run python ./bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全加固
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full

# uv 需要写入这些目录
ReadWritePaths=/opt/HIKARI_BOT/data /opt/HIKARI_BOT/logs /opt/HIKARI_BOT/.venv

[Install]
WantedBy=multi-user.target
